# files/etc/crontabs/root
#!/bin/sh

# =========================================================
# 1. Password root & Izin File
# =========================================================
echo "Setting root password..."
echo -e "root\nroot" | passwd root > /dev/null 2>&1

# Loop izin eksekusi untuk semua file .sh
for file in /sbin/free.sh /etc/rc.local /etc/custom_service/start_service.sh /etc/custom_service/free.sh /www/vnstati/vnstati.sh; do
    [ -f "$file" ] && chmod +x "$file"
done

# =========================================================
# 2. Konfigurasi UCI (System, Network, Firewall, uHTTPd)
# =========================================================
echo "Configuring UCI settings..."
uci batch <<EOF
    set system.@system[0].hostname='FRDM-X'
    set system.@system[0].timezone='WIB-7'
    set system.@system[0].zonename='Asia/Jakarta'
    
    delete system.ntp.server
    add_list system.ntp.server='0.id.pool.ntp.org'
    add_list system.ntp.server='1.id.pool.ntp.org'
    add_list system.ntp.server='://google.com'
    
    set network.mm=interface
    set network.mm.proto='modemmanager'
    set network.mm.device='auto'
    set network.mm.apn='internet'
    set network.mm.peerdns='1'
    set network.mm.metric='10'
    
    set network.wan=interface
    set network.wan.proto='dhcp'
    set network.wan.device='eth1'
    
    set network.tethering=interface
    set network.tethering.proto='dhcp'
    set network.tethering.device='usb0'
    
    delete network.wan6
    
    del_list firewall.@zone[1].network='wan6'
    add_list firewall.@zone[1].network='tethering'
    add_list firewall.@zone[1].network='wan'
    add_list firewall.@zone[1].network='mm'
    
    set uhttpd.main.index_page='index.php'
    set uhttpd.main.interpreter='.php=/usr/bin/php-cgi'
    
    commit network
    commit system
    commit firewall
    commit uhttpd
    commit dhcp
EOF

# =========================================================
# 3. Optimasi & Fixes
# =========================================================
echo "Applying optimizations..."
sed -i 's/^option check_signature/# option check_signature/g' /etc/opkg.conf 2>/dev/null

for dev in /sys/bus/usb/devices/*/power/control; do 
    [ -e "$dev" ] && echo "on" > "$dev" 2>/dev/null
done

[ -f /etc/php.ini ] && sed -i 's/^memory_limit.*/memory_limit = 64M/' /etc/php.ini

if [ -f /etc/config/ttyd ] && [ -f /bin/bash ]; then
    uci set ttyd.@ttyd[0].command='/bin/bash --login'
    uci commit ttyd
    /etc/init.d/ttyd enable
fi

# =========================================================
# 4. TinyFM, Nikki & OpenClash Symlinks
# =========================================================
mkdir -p /www/tinyfm
ln -sf / /www/tinyfm/rootfs
ln -sf /etc/openclash/ /www/tinyfm/openclash
ln -sf /etc/openclash/ /www/tinyfm/nikki

# =========================================================
# 5. Konfigurasi Jadwal (Cron Job / Schedule) Universal Port
# =========================================================
echo "Updating System Schedules (Universal Port)..."

cat > /etc/crontabs/root <<EOF
# Script by FRDMX

# auto renew ip lease (Mon-Sat 03:30 AM)
#30 3 * * 1,2,3,4,5,6 for dev in /dev/ttyUSB*; do [ -e "\$dev" ] && echo AT+CFUN=4 | atinout - "\$dev" -; done && sleep 3 && ifdown wan && sleep 3 && for dev in /dev/ttyUSB*; do [ -e "\$dev" ] && echo AT+CFUN=1 | atinout - "\$dev" -; done && sleep 3 && ifup wan

# auto restart modem (Sunday 03:30 AM)
#30 3 * * 0 for dev in /dev/ttyUSB*; do [ -e "\$dev" ] && echo AT^RESET | atinout - "\$dev" -; done && sleep 20 && ifdown wan && ifup wan

0 * * * * /sbin/free.sh >/dev/null 2>&1
*/5 * * * * /www/vnstati/vnstati.sh >/dev/null 2>&1
0 * * * * /bin/sh /etc/custom_service/start_service.sh
EOF

/etc/init.d/cron enable
/etc/init.d/cron restart

# =========================================================
# 6. Restart Services & Eksekusi Akhir
# =========================================================
echo "Restarting services..."
sysctl -p /etc/sysctl.conf >/dev/null 2>&1
/etc/init.d/network restart
/etc/init.d/uhttpd restart
/etc/init.d/ttyd restart

sleep 2
/etc/init.d/modemmanager restart

[ -f "/sbin/free.sh" ] && sh /sbin/free.sh
[ -f "/etc/custom_service/free.sh" ] && sh /etc/custom_service/free.sh

echo "Done! Router FRDM-X is ready."
exit 0
