# files/etc/rc.local
# Put your custom commands here that should be executed once
#sleep 5 && /usr/bin/x-gpio -r
#sleep 5 && /sbin/free.sh
#/sbin/jam bug.com
exit 0